# yarr/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant  Add-on: yarr
# Runs the Nginx daemon
# ==============================================================================

# Wait for yarr to become available
bashio::net.wait_for 7070

bashio::log.info "Starting NGinx..."

exec nginx
