# yarr/rootfs/etc/services.d/yarr/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: yarr
# Runs yarr
# ==============================================================================

bashio::log.info 'Starting yarr...'

# Run yarr
exec /usr/local/bin/yarr.sh
