# yarr/rootfs/etc/s6-overlay/s6-rc.d/yarr/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant App: yarr
# Runs yarr
# ==============================================================================

bashio::log.info 'Starting yarr...'

# Run yarr
exec /usr/local/bin/yarr.sh
